# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2ctl]
command = /usr/sbin/apache2ctl
user = root
autostart = true

